# rootfs/etc/services.d/turtlecoind/run
#!/usr/bin/with-contenv bash

WALLET_ADDRESS="${WALLET_ADDRESS:-}"
FEE_AMOUNT="${NODE_FEE:-0}"

echo "### TurtleCoin Node ###"
echo "Fee: ${FEE_AMOUNT} shells"
echo "Wallet address: ${WALLET_ADDRESS}"

s6-setuidgid "${TRTL_USER}" "${TRTL_PATH}/TurtleCoind" \
  --rpc-bind-ip "0.0.0.0" \
  --data-dir "${TRTL_HOME}/blockchain" \
  --db-threads 5 \
  --log-file "${TRTL_HOME}/TurtleCoind.log" \
  --load-checkpoints "${TRTL_HOME}/checkpoints.csv" \
  --fee-amount "${FEE_AMOUNT}" \
  --fee-address "${WALLET_ADDRESS}"
